# systemd/dhcp-admin-gui.service
[Unit]
Description=DHCP Web Admin GUI
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dhcp_server/services
EnvironmentFile=/opt/dhcp_server/services/.env
ExecStart=/opt/dhcp_server/.venv/bin/uvicorn admin_gui:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
